# brightness_scheduler.timer
[Unit]
Description=Adjust brightness every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
